# templates/etc/systemd/system/rstudio-pm.service.d/user.conf
[Service]
User={{ rstudio_package_manager_run_as_user }}
Group={{ rstudio_package_manager_run_as_group }}
EnvironmentFile=/etc/environment
