# templates/etc/systemd/system/rstudio-pm.service.d/user.conf
[Service]
User=rstudio-pm
Group=rstudio-pm
EnvironmentFile=/etc/environment
